# installer/linux/deb/resources/radikeep.service
[Unit]
Description=RadiKeep Web Application
After=network-online.target
Wants=network-online.target

[Service]
# 専用ユーザーで実行し、書き込み先は ReadWritePaths に限定する。
Type=simple
User=radikeep
Group=radikeep
WorkingDirectory=/opt/radikeep
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=RADIKEEP_BIND_ADDRESS=0.0.0.0
Environment=RADIKEEP_PORT=8085
Environment=ASPNETCORE_URLS=http://0.0.0.0:8085
Environment=RadiKeep__DbDirectory=/var/lib/radikeep/db
Environment=RadiKeep__LogDirectory=/var/log/radikeep
Environment=RadiKeep__DataProtectionKeysPath=/var/lib/radikeep/keys
Environment=RadiKeep__RecordFileSaveFolder=/var/lib/radikeep/record
Environment=RadiKeep__TemporaryFileSaveFolder=/var/lib/radikeep/temp
# /etc/default/radikeep でポートや保存先を上書き可能。
EnvironmentFile=-/etc/default/radikeep
ExecStart=/opt/radikeep/start-radikeep.sh
Restart=always
RestartSec=5
# systemd hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/radikeep /var/log/radikeep /etc/radikeep
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
